# clients/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
#ExecStart=-/sbin/agetty --autologin dsr --noclear %I $TERM
