# configuration/monitoring/crontab
5,15,25,25,45,55 * * * *   elgaard  /data/roprotokol/configuration/monitoring/bin/checkserver.sh
